# backup-server/crontab
# Run backup every 5 minutes for testing
*/5 * * * * /backup.sh >> /logs/cron.log 2>&1
