# backup-server/crontab
# Run backup every minute for testing
* * * * * /backup.sh >> /logs/cron.log 2>&1
